# systemd/system/abstracts-http-proxy.service
# Forwards connections arriving on port 80 (held by abstracts-http.socket) to the
# rootless Podman proxy container (nginx or Caddy) listening on 127.0.0.1:8080.
#
# systemd-socket-proxyd is part of systemd itself (no additional packages needed).
# It transparently proxies the raw TCP stream, so the proxy container sees the
# original client connection.
#
# See abstracts-http.socket for installation instructions.

[Unit]
Description=Forward HTTP (port 80) to Abstracts Explorer proxy container
Documentation=https://abstracts-explorer.readthedocs.io/en/latest/docker.html
Requires=abstracts-http.socket
After=abstracts-http.socket

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd 127.0.0.1:8080
# Allow time for the user proxy container to start after boot.
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
