# systemd/system/abstracts-https.socket
# System-level socket that opens port 443 on behalf of the rootless Podman stack.
# See abstracts-http.socket for installation instructions.

[Unit]
Description=HTTPS socket for Abstracts Explorer (port 443)
Documentation=https://abstracts-explorer.readthedocs.io/en/latest/docker.html

[Socket]
ListenStream=0.0.0.0:443
Accept=no

[Install]
WantedBy=sockets.target
