# systemd/abstracts-explorer.service
# Abstracts Explorer — systemd service unit (Docker Compose, existing-certificate variant)
#
# This unit starts the Abstracts Explorer stack defined in docker-compose.yml,
# which uses nginx with your own SSL certificate files (certs/cert.pem and
# certs/key.pem).  If you want automatic Let's Encrypt certificates, use
# abstracts-explorer-letsencrypt.service instead.
#
# QUICK START
# ===========
# 1. Place docker-compose.yml, .env, and the nginx/ and certs/ directories in
#    the deployment directory (default: /opt/abstracts-explorer).
#    Adjust the WorkingDirectory path below if you chose a different location.
#
# 2. Install and start the service:
#
#      sudo cp abstracts-explorer.service /etc/systemd/system/
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now abstracts-explorer
#
# 3. Check service status:
#
#      sudo systemctl status abstracts-explorer
#      sudo journalctl -u abstracts-explorer -f
#
# PODMAN USERS
# ============
# Replace "docker compose" with "podman-compose" (or "podman compose" on
# Podman 4.x+) in the ExecStart/ExecStop lines below, and change the
# docker.socket dependency to "After=network-online.target" only (Podman is
# daemonless).  For a rootless setup, install the unit as a user service
# instead (see abstracts-explorer-podman.service).

[Unit]
Description=Abstracts Explorer (Docker Compose — existing certificate)
Documentation=https://abstracts-explorer.readthedocs.io/en/latest/docker.html
# Docker socket must be available before this service starts.
Requires=docker.socket
After=docker.socket network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Path to the directory that contains docker-compose.yml and the .env file.
# Change this if you placed the files elsewhere.
WorkingDirectory=/opt/abstracts-explorer

# Pull the latest images before starting (comment out to skip during restart).
ExecStartPre=-/usr/bin/docker compose pull --quiet

# Start all containers in detached mode, remove orphaned containers.
ExecStart=/usr/bin/docker compose up -d --remove-orphans

# Stop all containers cleanly when the service is stopped or the system shuts down.
ExecStop=/usr/bin/docker compose down

# Allow enough time for images to pull and all containers to become healthy.
TimeoutStartSec=300
TimeoutStopSec=120

# Automatically restart the service if it exits unexpectedly.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
