# systemd/abstracts-explorer-letsencrypt.service
# Abstracts Explorer — systemd service unit (Docker Compose, Let's Encrypt variant)
#
# This unit starts the Abstracts Explorer stack defined in
# docker-compose.letsencrypt.yml, which obtains and renews TLS certificates
# automatically from Let's Encrypt.  If you already have your own SSL
# certificate files, use abstracts-explorer.service instead.
#
# PREREQUISITES
# =============
# Before starting this service for the first time you must obtain the initial
# certificate with Certbot in standalone mode (nginx must NOT be running yet):
#
#   docker run --rm \
#     -p 80:80 \
#     -v letsencrypt-certs:/etc/letsencrypt \
#     certbot/certbot certonly --standalone \
#     --domain abstracts.example.com \
#     --email your@email.com \
#     --agree-tos --non-interactive
#
# Replace "abstracts.example.com" with your real domain everywhere in
# docker-compose.letsencrypt.yml and nginx/nginx.letsencrypt.conf first.
#
# QUICK START
# ===========
# 1. Place docker-compose.letsencrypt.yml, .env, and the nginx/ directory in
#    the deployment directory (default: /opt/abstracts-explorer).
#    Adjust the WorkingDirectory path below if you chose a different location.
#
# 2. Obtain the initial certificate (see PREREQUISITES above).
#
# 3. Install and start the service:
#
#      sudo cp abstracts-explorer-letsencrypt.service /etc/systemd/system/
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now abstracts-explorer-letsencrypt
#
# 4. Check service status:
#
#      sudo systemctl status abstracts-explorer-letsencrypt
#      sudo journalctl -u abstracts-explorer-letsencrypt -f
#
# AUTOMATIC CERTIFICATE RENEWAL
# ==============================
# The certbot container (defined in docker-compose.letsencrypt.yml) checks for
# renewal every 12 hours and renews the certificate when fewer than 30 days
# remain.  After renewal, nginx must be reloaded to pick up the new certificate.
# Add a daily cron job to do this automatically:
#
#   sudo crontab -e
#   # Add the following line (reloads nginx daily at 3 AM to pick up renewed certificates):
#   0 3 * * * /usr/bin/docker exec abstracts-nginx nginx -s reload
#
# PODMAN USERS
# ============
# Replace "docker compose -f docker-compose.letsencrypt.yml" with
# "podman-compose -f docker-compose.letsencrypt.yml" in the ExecStart/ExecStop
# lines, and change "docker exec" to "podman exec" in the cron job above.
# Remove the docker.socket dependency and use "After=network-online.target"
# only.  For a rootless setup, install as a user service instead
# (see abstracts-explorer-letsencrypt-podman.service).

[Unit]
Description=Abstracts Explorer (Docker Compose — Let's Encrypt)
Documentation=https://abstracts-explorer.readthedocs.io/en/latest/docker.html
# Docker socket must be available before this service starts.
Requires=docker.socket
After=docker.socket network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Path to the directory that contains docker-compose.letsencrypt.yml and .env.
# Change this if you placed the files elsewhere.
WorkingDirectory=/opt/abstracts-explorer

# Pull the latest images before starting (comment out to skip during restart).
ExecStartPre=-/usr/bin/docker compose -f docker-compose.letsencrypt.yml pull --quiet

# Start all containers in detached mode, remove orphaned containers.
ExecStart=/usr/bin/docker compose -f docker-compose.letsencrypt.yml up -d --remove-orphans

# Stop all containers cleanly when the service is stopped or the system shuts down.
ExecStop=/usr/bin/docker compose -f docker-compose.letsencrypt.yml down

# Allow enough time for images to pull and all containers to become healthy.
TimeoutStartSec=300
TimeoutStopSec=120

# Automatically restart the service if it exits unexpectedly.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
