# systemd/abstracts-explorer.service
# Abstracts Explorer — systemd service unit (Docker Compose)
#
# This unit starts the Abstracts Explorer stack defined in docker-compose.yml.
# It works with any of the provided Compose configurations (existing certificate
# or Let's Encrypt); copy the chosen file to this directory as docker-compose.yml.
#
# SETUP
# =====
# 1. Choose the Compose file that best matches your setup:
#      - docker-compose.yml              (existing SSL certificate)
#      - docker-compose.letsencrypt.yml  (automatic Let's Encrypt certificate)
#
# 2. Copy the chosen file and its nginx configuration to the deployment directory
#    (default: /opt/abstracts-explorer), renaming the Compose file to
#    docker-compose.yml if necessary:
#
#      sudo mkdir -p /opt/abstracts-explorer
#      sudo cp docker-compose.yml /opt/abstracts-explorer/
#      sudo cp -r nginx/ /opt/abstracts-explorer/
#      # If using the existing-certificate variant also copy the certs/ directory:
#      sudo cp -r certs/ /opt/abstracts-explorer/
#
# 3. Add a .env file with any required tokens (e.g. LLM_BACKEND_AUTH_TOKEN):
#
#      sudo cp .env /opt/abstracts-explorer/
#
# 4. Edit /opt/abstracts-explorer/docker-compose.yml to set your configuration
#    (domain name, model names, passwords, etc.).
#
# 5. If using Let's Encrypt, obtain the initial certificate first (nginx must
#    not be running yet):
#
#      docker run --rm -p 80:80 \
#        -v letsencrypt-certs:/etc/letsencrypt \
#        certbot/certbot certonly --standalone \
#        --domain your.domain.com \
#        --email your@email.com \
#        --agree-tos --non-interactive
#
# 6. Install and start the service:
#
#      sudo cp abstracts-explorer.service /etc/systemd/system/
#      sudo systemctl daemon-reload
#      sudo systemctl enable --now abstracts-explorer
#
# 7. Check service status:
#
#      sudo systemctl status abstracts-explorer
#      sudo journalctl -u abstracts-explorer -f
#
# PODMAN USERS
# ============
# Use abstracts-explorer-podman.service instead of this file.

[Unit]
Description=Abstracts Explorer (Docker Compose)
Documentation=https://abstracts-explorer.readthedocs.io/en/latest/docker.html
# Docker socket must be available before this service starts.
Requires=docker.socket
After=docker.socket network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Path to the directory that contains docker-compose.yml and the .env file.
# Change this if you placed the files elsewhere.
WorkingDirectory=/opt/abstracts-explorer

# Pull the latest images before starting (comment out to skip during restart).
ExecStartPre=-/usr/bin/docker compose pull --quiet

# Start all containers in detached mode, remove orphaned containers.
ExecStart=/usr/bin/docker compose up -d --remove-orphans

# Stop all containers cleanly when the service is stopped or the system shuts down.
ExecStop=/usr/bin/docker compose down

# Allow enough time for images to pull and all containers to become healthy.
TimeoutStartSec=300
TimeoutStopSec=120

# Automatically restart the service if it exits unexpectedly.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
