# systemd/user/caddy/abstracts-postgres.container
[Unit]
Description=PostgreSQL database for Abstracts Explorer (Caddy variant)
After=abstracts-postgres-data-volume.service
Requires=abstracts-postgres-data-volume.service

[Container]
Image=docker.io/postgres:16-alpine
ContainerName=abstracts-postgres
Network=abstracts.network

Volume=abstracts-postgres-data.volume:/var/lib/postgresql/data

Environment=POSTGRES_DB=abstracts
Environment=POSTGRES_USER=abstracts
# Store the database password as a Podman secret (see setup instructions in docs/docker.md).
#   podman secret create postgres-password -
# Then uncomment the Secret line and remove the Environment line below.
# Secret=postgres-password,type=env,target=POSTGRES_PASSWORD
Environment=POSTGRES_PASSWORD=change_me_in_production

HealthCmd=pg_isready -U abstracts
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=10s

[Service]
Restart=on-failure
TimeoutStartSec=120

[Install]
WantedBy=default.target
